# templates/claude-control@.service
[Unit]
Description=Claude Code Remote Control (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/claude-control/%i.env
WorkingDirectory=%%WORKING_DIR%%
ExecStart=/bin/bash %h/.config/claude-control/wrapper-%i.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
