# templates/com.claude-control.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-control.%%PROJECT_NAME%%</string>

    <key>ProgramArguments</key>
    <array>
        <string>%%CLAUDE_BIN%%</string>
        <string>remote-control</string>
        <string>--spawn</string>
        <string>worktree</string>
        <string>--capacity</string>
        <string>%%CAPACITY%%</string>
        <string>--name</string>
        <string>%%SESSION_NAME%%</string>
        <string>--no-create-session-in-dir</string>
    </array>

    <key>WorkingDirectory</key>
    <string>%%PROJECT_DIR%%</string>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>%%PATH%%</string>
    </dict>

    <key>StandardOutPath</key>
    <string>%%LOG_DIR%%/%%PROJECT_NAME%%.log</string>

    <key>StandardErrorPath</key>
    <string>%%LOG_DIR%%/%%PROJECT_NAME%%.log</string>
</dict>
</plist>
